# buster/profile.d/10_meta_compatible.sh
#CELL=`cat /etc/openafs/ThisCell`
CELL=ics.muni.cz
export CELL

# kompatibilita s SuSE ~/.profile
PROFILEREAD=yes
